# issue-fancy.crontab
PATH=/sbin:/bin:/usr/sbin:/usr/bin

00,15,30,45 * * * *	root	/sbin/issue-fancy-gen
